# supervisor.ini
[program:laravel-octane]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan octane:start --watch --host=0.0.0.0 --port=80
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
